# example/vim.conf
[program:vim]
command=xterm -display :0 -e /bin/bash -l -c "vim"
stdout_logfile=/tmp/vim.log
stderr_logfile=/tmp/vim.log
user=user
environment=DISPLAY=":0"
autorestart=true
priority=999
